# payload/Library/LaunchDaemons/com.github.sphen13.bzinstall.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
	<key>Label</key>
	<string>com.github.sphen13.bzinstall</string>
	<key>ProgramArguments</key>
	<array>
		<string>/Library/Scripts/bz/backblaze-deployer.sh</string>
	</array>
	<key>RunAtLoad</key>
	<false/>
	<key>KeepAlive</key>
	<dict>
		<key>PathState</key>
		<dict>
			<key>/Library/Scripts/bz/.trigger</key>
			<true/>
		</dict>
	</dict>
	<key>OnDemand</key>
	<true/>
</dict>
</plist>
